# supervisor/profile/redis-server.ini
[program:redis-server]
command=redis-server
directory=/home/Projects/redis-server
autorestart=true
startsecs=3
startretries=3
redirect_stderr=true
stdout_logfile=/tmp/log/redis-server.out.log
stderr_logfile=/tmp/log/redis-server.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=root
priority=999
numprocs=1
process_name=%(program_name)s_%(process_num)02d
environment=HOME=/root